# deploy/systemd/neuramd-backup.timer
[Unit]
Description=Run NeuraMD pg_dump daily at 03:00
Requires=neuramd-backup.service

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
