# deploy/systemd/neuramd-worker.service
[Unit]
Description=NeuraMD Solid Queue worker on AIrch
After=network-online.target postgresql.service neuramd-web.service
Wants=network-online.target
Requires=postgresql.service

[Service]
Type=simple
User=rafael
Group=rafael
WorkingDirectory=/home/rafael/apps/NeuraMD
Environment=RAILS_ENV=production
Environment=RAILS_LOG_TO_STDOUT=1
Environment=POSTGRES_USER=postgres
Environment=POSTGRES_HOST=127.0.0.1
Environment=POSTGRES_PORT=5432
Environment=NEURAMD_DATABASE=neura_md_production
Environment=NEURAMD_CACHE_DATABASE=neura_md_production_cache
Environment=NEURAMD_QUEUE_DATABASE=neura_md_production_queue
Environment=NEURAMD_LEGACY_DATABASE=neura_md_legacy
Environment=NEURAMD_TENTACLES_ENABLED=1
Environment=PATH=/home/rafael/.local/share/mise/shims:/home/rafael/.local/bin:/usr/local/bin:/usr/bin
ExecStart=/home/rafael/.local/share/mise/shims/bundle exec bin/jobs
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=neuramd-worker

# Tentacle cwd whitelist. The cron tick job boots tentacles from this
# worker, so any custom cwd root must be mirrored in ReadWritePaths
# (specifically the `<root>/tmp/tentacles` subpath) or the worktree
# provisioning call will fail under ProtectHome=read-only. Leaving it
# unset falls back to Tentacles::BootConfig::DEFAULT_CWD_ALLOWED_PREFIXES.
# Environment=NEURAMD_TENTACLE_CWD_ROOTS=/home/rafael/projects/

# Shared dtach runtime directory. The web unit owns this directory
# (RuntimeDirectoryPreserve=yes), but the supervisor job running here
# needs write access to sweep orphaned sockets. Declaring the same
# RuntimeDirectory entry mounts the path with matching ownership.
RuntimeDirectory=nm-tentacles
RuntimeDirectoryMode=0750
RuntimeDirectoryPreserve=yes

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/home/rafael/apps/NeuraMD/tmp /home/rafael/apps/NeuraMD/log /home/rafael/apps/NeuraMD/storage /home/rafael/apps/NeuraMD/.git /home/rafael/.claude /home/rafael/.claude.json

[Install]
WantedBy=multi-user.target
