# deploy/systemd/neuramd-worker.service.d/git-mirrors.conf
[Service]
# Mirror of the web unit's git-mirrors.conf. WorktreeService /
# tentacle push paths may execute from either process, so both units
# need rw access to the bare mirror. See `83228dc2` for the full
# diagnosis from the DevOps tentacle.
ReadWritePaths=/home/rafael/git-mirrors
