# deploy/systemd/neuramd-worker.service.d/npm.conf
[Service]
# Parallels the web drop-in — the supervisor/cron paths may eventually
# need to shell out to npx too, and keeping the two units symmetric
# avoids divergence when a job is moved between them.
ReadWritePaths=/home/rafael/.npm
