# deploy/systemd/neuramd-worker.service.d/tentacles-git.conf
[Service]
ReadWritePaths=/home/rafael/apps/NeuraMD/.git
