# deploy/systemd/neuramd-worker.service.d/tentacles.conf
[Service]
Environment=NEURAMD_TENTACLES_ENABLED=1

# Shared dtach runtime directory. The web unit owns this directory
# (RuntimeDirectoryPreserve=yes), but the supervisor job running here
# needs write access to sweep orphaned sockets. Declaring the same
# RuntimeDirectory entry mounts the path with matching ownership.
RuntimeDirectory=nm-tentacles
RuntimeDirectoryMode=0750
RuntimeDirectoryPreserve=yes
