# deploy/systemd/neuramd-worker.service.d/trusted-workspaces.conf
[Service]
# Mirrors the web unit: cron tentacles spawning through this worker
# need the same allowlist so their WorktreeService.ensure calls link
# master.key on trusted workspaces.
Environment=NEURAMD_TRUSTED_WORKSPACE_KEYS=/home/rafael/workspaces/neuramd
