# deploy/systemd/neuramd-worker.service.d/workspace.conf
[Service]
ReadWritePaths=/home/rafael/.codex
ReadWritePaths=/home/rafael/workspaces
