# deploy/systemd/neuramd-web.service.d/claude-home.conf
[Service]
ReadWritePaths=/home/rafael/.claude
ReadWritePaths=/home/rafael/.claude.json
