# deploy/systemd/neuramd-web.service.d/git-mirrors.conf
[Service]
# Tentacle worktrees push to /home/rafael/git-mirrors/NeuraMD.git (the
# bare mirror that the post-receive hook auto-deploys from). Without
# this path in the rw allowlist, `git push` from inside a tentacle
# fails with "unable to create temporary object directory" — the
# unpacker can't write to objects/incoming-* and the receive-pack
# rejects the push. Documented in note `83228dc2` and reported by the
# DevOps tentacle 2026-04-26 after blocking Tentacles::IdentityFile push.
ReadWritePaths=/home/rafael/git-mirrors
