# deploy/systemd/neuramd-web.service.d/npm.conf
[Service]
# Agents running `claude` inside a tentacle spawn MCP servers as
# subprocesses; npx writes to ~/.npm/_npx when resolving packages like
# chrome-devtools-mcp. Without this path in the whitelist, npx fails
# with EROFS before the MCP even starts. Scope: user-owned cache dir,
# safe to expose to the same sandboxed process.
ReadWritePaths=/home/rafael/.npm
