# deploy/systemd/neuramd-web.service.d/tentacles.conf
[Service]
Environment=NEURAMD_TENTACLES_ENABLED=1

# Dtach runtime directory. Detached tentacle sessions store their Unix
# socket and pidfile under /run/nm-tentacles. RuntimeDirectoryPreserve
# keeps the directory across `systemctl restart` so the next boot's
# bootstrap_sessions! can reattach. systemd wipes it on reboot (tmpfs),
# which is fine — processes die on reboot anyway.
RuntimeDirectory=nm-tentacles
RuntimeDirectoryMode=0750
RuntimeDirectoryPreserve=yes
