# deploy/systemd/neuramd-web.service.d/tentacles-git.conf
[Service]
# Tentacle worktrees under tmp/tentacles/<uuid> share the parent repo's
# .git directory — `git worktree add` writes HEAD/config/logs/refs there,
# so the sandbox must whitelist it. Mirrored in the worker unit via its
# own drop-in because WorktreeService may run from either process.
ReadWritePaths=/home/rafael/apps/NeuraMD/.git
