# deploy/systemd/neuramd-web.service.d/trusted-workspaces.conf
[Service]
# Explicit allowlist of workspace paths that may receive the host
# app's config/master.key via WorktreeService.ensure. Empty/unset env
# keeps the feature off; any workspace listed here is trusted to boot
# Rails with the running neuramd-web's encryption key.
#
# Admin responsibility: only list workspaces whose contents (repo
# HEAD + working tree) are under the same operational trust as the
# running app. A fresh clone of a public remote placed under
# /home/rafael/workspaces/ is NOT automatically trusted — it must be
# explicitly added here.
Environment=NEURAMD_TRUSTED_WORKSPACE_KEYS=/home/rafael/workspaces/neuramd
