# deploy/systemd/neuramd-web.service.d/workspace.conf
[Service]
# Added with tentacle_workspace feature (PR #13). /home/rafael/.codex for
# codex CLI state; /home/rafael/workspaces for Tentacles::Workspace root.
ReadWritePaths=/home/rafael/.codex
ReadWritePaths=/home/rafael/workspaces
